# etc/cron.d/lipsync
# have lipsync sync with the server every 30 minutes to catch changes from other nodes
09,39 *     * * *     $username	/usr/bin/lipsync -ltd --delete --rsh="ssh -l $username -p $port" $lipsync_dir_local/ $remote_server:$lipsync_dir_remote/
